# goblin/systemd/goblin.service
[Unit]
Description=MonsterBox Goblin Media Player
Documentation=https://github.com/arwpc/MonsterBox
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
User=remote
WorkingDirectory=/home/remote/goblin
Environment="GOBLIN_ID=goblin-one"
Environment="GOBLIN_PORT=3001"
Environment="NODE_ENV=production"
Environment="MONSTERBOX_URL=http://192.168.8.140:3000"

# Pre-start cleanup and display claim (ignore failures)
ExecStartPre=-/usr/bin/pkill -f ffplay
ExecStartPre=-/usr/bin/pkill -f mpv
ExecStartPre=-/usr/bin/pkill -f vlc
ExecStartPre=-/bin/bash -lc 'lsof -t -i :3001 | xargs -r kill -9'
ExecStartPre=-/usr/bin/test -x /usr/local/bin/claim_display.sh && /usr/local/bin/claim_display.sh

# Start server
ExecStart=/usr/bin/node /home/remote/goblin/src/server.js

# Restart policy - production hardened
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
KillMode=mixed

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=append:/home/remote/goblin/logs/goblin.log
StandardError=append:/home/remote/goblin/logs/goblin-error.log
SyslogIdentifier=goblin

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Watchdog - restart if service becomes unresponsive
WatchdogSec=60

[Install]
WantedBy=multi-user.target

